# com.user.claude-usage-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-usage-tracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/pc00002270/.local/bin/claude-usage-tracker</string>
        <string>daemon</string>
        <string>--interval</string>
        <string>900</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/claude-usage-tracker.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-usage-tracker.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>
</dict>
</plist>
